# scripts/bin/membraned.service
[Unit]
Description=Membrane Daemon
After=network.target

[Service]
ExecStart=/usr/share/membrane/bin/membraned
Restart=on-failure

[Install]
WantedBy=multi-user.target